# screener-sml2.service
[Unit]
Description=StockTipsBot — SML2 Screener (WebSocket, $0.50–$5.00 paper)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=defectuous
WorkingDirectory=/home/defectuous/StockTipsBot
ExecStart=/home/defectuous/StockTipsBot/.venv/bin/python run_sml2_screener.py
Restart=on-failure
RestartSec=30
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
