# stockbot.service
[Unit]
Description=StockTipsBot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Change User and WorkingDirectory to match your Pi setup
User=defectuous
WorkingDirectory=/home/defectuous/StockTipsBot
ExecStart=/home/defectuous/StockTipsBot/.venv/bin/python main.py
Restart=on-failure
RestartSec=30
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
